# supervisord.conf
[supervisord]
loglevel=debug

[program:openvpn]
command=openvpn --config /root/.ssh/client.ovpn

